# scripts/backup.service
[Unit]
Description=UP Schedule Generator Backup Service
Documentation=https://github.com/yourusername/schedgen
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/schedgen
EnvironmentFile=/opt/schedgen/.env
ExecStart=/opt/schedgen/scripts/backup-all.sh
StandardOutput=append:/var/log/schedgen-backup.log
StandardError=append:/var/log/schedgen-backup.log

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/schedgen/backups /var/log

[Install]
WantedBy=multi-user.target
